# theia/deployment/pluribus-theia.service
[Unit]
Description=Pluribus Theia Vision Agent
After=network.target pluribus-omega.service
Requires=pluribus-omega.service

[Service]
Type=simple
User=root
WorkingDirectory=/pluribus
Environment=PYTHONPATH=/pluribus
Environment=PLURIBUS_VPS_MODE=1
ExecStart=/usr/bin/python3 -m theia.deployment.run_daemon
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
